# radar_parlamentar/cron/crontab
# crontab -e to edit crontab
# crontab syntax:
# mm hh dd MM ss script
# mm = minuto(0-59)
# hh = hora(0-23)
# dd = dia(1-31)
# MM = mes(1-12)
# ss = dia_da_semana(0-6)

0 4 * * 0 dump-radar.sh

